# proccesslogger.service
[Unit]
Description=processlogger

[Service]
ExecStart=/usr/bin/processlogger
Restart=always
RestartSec=12

[Install]
WantedBy=default.target