# etc/processlogger.service
[Unit]
Description=proccesslogger

[Service]
ExecStart=#idk
Type=simple